# supervisor_mongo_vazio.ini
[program:ajna_kerasvazio_deploy]
environment=MONGODB_URI=mongodb://127.0.0.1
directory=/home/ivan/ajna_tfod_deploy
command=/home/ivan/ajna_tfod_deploy/venv_tf23/bin/python mongo_periodic.py
autostart=true
autorestart=true
user=ivan
stderr_logfile=/var/log/ajna_kerasvazio_deploy/error.log
stdout_logfile=/var/log/ajna_kerasvazio_deploy/access.log