# supervisor_api.conf
[program:ajna_tfod_api]
directory=/home/ivan/ajna_tfod_deploy
command=/home/ivan/ajna_tfod_deploy/venv/bin/python app.py
autostart=true
autorestart=true
user=ivan
stderr_logfile=/var/log/ajna_tfod_deploy/error.log
stdout_logfile=/var/log/ajna_tfod_deploy/access.log