# supervisor_mongo.ini
[program:ajna_tfod_deploy]
environment=MONGODB_URI=mongodb://user:password@127.0.0.1/ajna
directory=/data/ajna_deploy/ajna_tfod_deploy
command=/data/ajna_deploy/ajna_tfod_deploy/venv_tf22/bin/python mongo_periodic.py
autostart=true
autorestart=true
user=ajna
stderr_logfile=/var/log/ajna_tfod_deploy/error.log
stdout_logfile=/var/log/ajna_tfod_deploy/access.log