# supervisor_mongo_motor_reefer.ini
[program:ajna_motor_reefer_deploy]
environment=MONGODB_URI=mongodb://127.0.0.1;SQL_URI=mysql+pymysql://mysql.labin.rf08.srf:3306/dbmercante
directory=/home/ivan/ajna_tfod_deploy
command=/home/ivan/ajna_tfod_deploy/venv_torch/bin/python motor_reefer/mongo_periodic.py
autostart=true
autorestart=true
user=ivan
stderr_logfile=/var/log/ajna_motor_reefer_deploy/error.log
stdout_logfile=/var/log/ajna_motor_reefer_deploy/access.log